# etc/systemd/system/low-battery-check.service
[Unit]
Description=Low battery level check and notification

[Service]
Type=oneshot
ExecStart=exec /etc/profile.d/low-battery-check.sh
